# examples/systemd/ops-light-secrets-server.service
[Unit]
Description=ops-light-secrets-server (data plane)
Documentation=file:docs/operating.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=olss
Group=olss
WorkingDirectory=/var/lib/olss
Environment=OLSS_DATA_DIRECTORY=/var/lib/olss
Environment=OLSS_CONFIG=/etc/olss/config.toml
ExecStart=/usr/local/bin/ops-light-secrets-server serve
Restart=on-failure
RestartSec=2
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/olss
CapabilityBoundingSet=
AmbientCapabilities=
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictRealtime=yes
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
