# examples/systemd/ops-light-secrets-server-control.socket
[Unit]
Description=ops-light-secrets-server owner control socket
Documentation=file:docs/operating.md

[Socket]
ListenStream=/run/olss/control.sock
SocketMode=0600
SocketUser=olss
SocketGroup=olss
DirectoryMode=0750

[Install]
WantedBy=sockets.target
